# Week-05-Scraper-API-Service/Project-Scraper-API-Service/debian/scraper-api-service/DEBIAN/prerm
#!/bin/bash
set -e

systemctl stop scraper-api.service || true
systemctl disable scraper-api.service || true

exit 0
